# systemd/user/screentracker.service
[Unit]
Description=Screen activity tracker (auto-detect X11)
After=graphical-session.target plasma-workspace.target

[Service]
Type=simple
ExecStart=/bin/bash -c '\
  DISPLAY=${DISPLAY:-:0}; \
  XAUTH=$(find /tmp -maxdepth 1 -type f -user $USER -name "xauth_*" | head -n1); \
  export DISPLAY XAUTHORITY=${XAUTH:-$XAUTHORITY}; \
  echo "Using DISPLAY=$DISPLAY, XAUTHORITY=$XAUTHORITY"; \
  exec "$HOME/.local/bin/screentracker.py" \
'
Restart=on-failure

[Install]
WantedBy=default.target
